# scripts/fastdds_discovery_server.service
[Unit]
Description=FastDDS Discovery Server (custom_turtlebot_nav_ws)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
ExecStart=/bin/bash -c 'source /opt/ros/humble/setup.bash && exec fast-discovery-server -i 0 --listening-addresses 0.0.0.0:11811'
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
